# buildroot/buildroot/board/amlogic/mesonc3_aw401/rootfs/etc/udev/rules.d/adb.rules
ACTION=="change", SUBSYSTEM=="android_usb", RUN+="/etc/rewrite_usb_udc.sh"
